# deploy/roles/precis_worker_agent/templates/precis-worker-agent.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
{# {{ ansible_managed }} #}
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.precis.worker-agent</string>

  {# bash wrapper makes /bin/bash the TCC "responsible process". Same
     pattern as precis_worker. #}
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>exec {{ precis_worker_agent_venv }}/bin/precis worker --profile agent --batch-size {{ precis_worker_agent_batch_size }} --idle-seconds {{ precis_worker_agent_idle_seconds }}</string>
  </array>

  {# Factory Phase-2 slice 2B: runs as `deploy` (was hermes). deploy owns
     the colima docker socket (0600) the §13 container executor needs, plus
     its own OAuth (~deploy/.claude_oauth_token, bootstrapped into
     CLAUDE_CODE_OAUTH_TOKEN by ensure_oauth_token) and its own
     ~/.claude/mcp.json (this role renders it from the shared asa_bot
     template). plan_tick uses --permission-mode bypassPermissions so no
     ~/.claude/settings.json dependency. dream stays its own hermes plist. #}
  <key>UserName</key>
  <string>deploy</string>

  <key>WorkingDirectory</key>
  <string>/Users/deploy</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/deploy</string>
    <key>PATH</key>
    <string>{{ precis_worker_agent_venv }}/bin:/opt/homebrew/bin:/Library/TeX/texbin:/usr/local/bin:/usr/bin:/bin</string>
    <key>PRECIS_DATABASE_URL</key>
    <string>{{ precis_worker_agent_dsn }}</string>
    <key>PRECIS_MCP_CONFIG</key>
    <string>{{ precis_worker_agent_mcp_config }}</string>
    <key>PRECIS_SOURCE</key>
    <string>precis-worker-agent</string>
    {# PRECIS_PROCESS tags every row in worker_logs so
       `precis logs --process precis-worker-agent` filters cleanly. #}
    <key>PRECIS_PROCESS</key>
    <string>precis-worker-agent</string>
    {# Gate envs — set to "1" to enable the corresponding review
       pass. Clearing any of these mutes the pass without
       uninstalling. #}
    <key>PRECIS_STRUCTURAL_REVIEW</key>
    <string>1</string>
    <key>PRECIS_DEEP_REVIEW</key>
    <string>1</string>
    {% if precis_worker_agent_load_ceiling %}
    <key>PRECIS_LOAD_CEILING</key>
    <string>{{ precis_worker_agent_load_ceiling }}</string>
    {% endif %}
    {# Planner-coroutine guardrails. Defaults are 10/2/20 in code;
       bumped here so the first real production planner can converge
       on a strategic with ~5-level decomposition before hitting the
       tick cap, and the daily ceiling allows ~10 planning sessions
       a day. Override at the host-vars level if a budget shift is
       needed. #}
    <key>PRECIS_MAX_TICKS</key>
    <string>{{ precis_max_ticks | default(10000) }}</string>
    <key>PRECIS_MAX_TODO_USD</key>
    <string>{{ precis_max_todo_usd | default(5.0) }}</string>
    <key>PRECIS_DAILY_COST_CEILING</key>
    <string>{{ precis_daily_cost_ceiling | default(50.0) }}</string>
    {# CPU-shares passthrough for containerized jobs this worker spawns
       (docker run --cpu-shares). Soft weight only — no core fencing on
       the Macs (PRECIS_JOB_CPUSET is set on the Linux/systemd node only). #}
    <key>PRECIS_JOB_CPU_SHARES</key>
    <string>{{ precis_job_cpu_shares | default('') }}</string>
    {# Autonomous quest loop (rung 4d, agent-profile only). The allocator
       pass mints frontier ticks + structure-relax compute; the weekly
       budget caps its spend (distinct from the global $/day breaker).
       Compute relax rides the existing the GPU node DFT lane.
       GATED OFF by default: the loop's frontier ticks are PAID calls, so
       it only does useful work when the global budget breaker
       (PRECIS_BUDGET_DAILY_USD) has headroom — otherwise every tick
       claims, hits the tripped breaker, fails, and re-claims in a
       ~0.5s spin. Flip `precis_quest_loop_enabled: true` in
       the gateway's host_vars once the daily cap is raised. #}
    <key>PRECIS_QUEST_LOOP_ENABLED</key>
    <string>{{ '1' if precis_quest_loop_enabled | default(false) else '0' }}</string>
    <key>PRECIS_QUEST_WEEKLY_BUDGET</key>
    <string>{{ precis_quest_weekly_budget | default(25) }}</string>
    {# Catalyst-quest barrier lane. dispatch_catpath runs HERE (quest ticks
       are claude_inproc, agent-profile only) and reads this env to pin the
       minted `catpath_explore` job at the GPU node; unset → it force-EMTs
       and the job never reaches the GPU node's MACE. Mirrors the precis-web plist,
       which sets the same var for the pathway-kind MCP surface. #}
    {% if (precis_capabilities.catpath | default([])) | length > 0 %}
    <key>PRECIS_CATPATH_ROUTE_NODE</key>
    <string>{{ precis_capabilities.catpath | first }}</string>
    {% endif %}
    {# diagram_propose autonomous drawer — agentic path (claude_inproc). #}
    <key>PRECIS_DIAGRAM_AGENTIC</key>
    <string>1</string>
    {# card_forge morning pass — `act` applies the retire/rework ladder
       (default `report` = observe-only). Composes here (claude_inproc). #}
    <key>PRECIS_CARD_FORGE_AUTONOMY</key>
    <string>{{ precis_card_forge_autonomy | default('act') }}</string>
    {# §13 container executor — CAPABLE but OFF (slice 2B). The docker CLI
       (homebrew) talks to deploy's colima VM where the precis-agent image
       is loaded (Track 1). Setting PRECIS_AGENT_CONTAINER=1 (slice 2C, via
       host_vars once plan_tick is verified under the deploy run-as) makes
       call_claude_agent passes (reviews/diagram) run in a throwaway
       container. plan_tick/fix_gripe use separate seams — unaffected. #}
    <key>PRECIS_CONTAINER_BIN</key>
    <string>{{ precis_worker_agent_container_bin }}</string>
    <key>PRECIS_AGENT_IMAGE</key>
    <string>{{ precis_agent_image }}</string>
    {% if precis_agent_container_enabled | default(false) %}
    <key>PRECIS_AGENT_CONTAINER</key>
    <string>1</string>
    {% endif %}
    {# Shared API keys — Perplexity, EPO OPS, patent raw cache.
       See inventory/group_vars/all/precis_env.yml. Same block in
       every precis-* daemon so capability parity is automatic. #}
    {% for env_key, env_val in precis_shared_env.items() %}
    <key>{{ env_key }}</key>
    <string>{{ env_val }}</string>
    {% endfor %}
  </dict>

  <key>RunAtLoad</key>
  <true/>

  {# Long-running daemon. KeepAlive=true so launchd restarts on
     crash; ThrottleInterval prevents a crash-loop from chewing
     CPU. Matches precis_worker. #}
  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>30</integer>

  {# Jetsam protection. This node co-hosts local llama.cpp inference
     servers that wire down the bulk of RAM (~70G+ resident, most of
     the machine's memory shows as wired). Under that pressure macOS
     terminates low-band background LaunchDaemons for "efficiency"
     (`launchctl print` reports `immediate reason = inefficient`) —
     and this worker, though tiny (~100 MB RSS), was being culled
     ~200x/day. Each kill that lands mid-tick orphans a claimed
     `plan_tick` (STATUS:running with no live runner) → the sweeper
     fails it `swept:claim-orphaned` → the parent todo child-fails and
     never converges. This host is the ONLY place plan_tick runs, so
     the worker must outlive the memory pressure. ProcessType
     Interactive raises its jetsam band toward the foreground tier so
     it is culled last, not first, and disables App-Nap-style timer
     throttling. The lightweight footprint makes the elevation cheap. #}
  <key>ProcessType</key>
  <string>Interactive</string>

  <key>StandardOutPath</key>
  <string>/var/log/precis-worker-agent.log</string>
  <key>StandardErrorPath</key>
  <string>/var/log/precis-worker-agent.log</string>
</dict>
</plist>
