# deploy/roles/llamacpp/templates/llamaswap.service.j2
# {{ ansible_managed }}
[Unit]
Description=llama-swap (llama.cpp multi-model proxy)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
Environment=HOME=/opt/llamacpp
Environment=PATH={{ llamacpp_bin_dir }}:/usr/local/bin:/usr/bin:/bin
ExecStart={{ llamacpp_bin_dir }}/llama-swap \
    -config /opt/llamacpp/etc/llama-swap.yaml \
    -listen 0.0.0.0:{{ llama_swap_port_linux }}
Restart=always
RestartSec=5
LimitNOFILE=65536
LimitMEMLOCK=infinity
# De-prioritize heavy compute so it can't starve interactive/system work
# (e.g. sshd). CPUAffinity is omitted on hosts that don't reserve cores.
Nice={{ precis_compute_nice | default(10) }}
{% if precis_compute_cpuaffinity | default('') %}CPUAffinity={{ precis_compute_cpuaffinity }}{% endif %}

[Install]
WantedBy=multi-user.target
