# deploy/roles/precis_reconcile/templates/precis-reconcile.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
{# {{ ansible_managed }} #}
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.precis.reconcile</string>

  <key>ProgramArguments</key>
  <array>
    <string>{{ precis_reconcile_venv }}/bin/precis</string>
    <string>reconcile-duplicates</string>
    <string>--apply</string>
  </array>

  <!-- Nightly one-shot. Collapses live paper refs sharing a pdf_sha256
       to the best survivor (DOI then non-junk title then most authors
       then lowest id). Cheap: one indexed GROUP BY, usually 0 groups. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>{{ precis_reconcile_hour }}</integer>
    <key>Minute</key><integer>{{ precis_reconcile_minute }}</integer>
  </dict>

  <!-- Transient timer, not a daemon: each fire runs one sweep and exits. -->
  <key>RunAtLoad</key>
  <false/>
  <key>KeepAlive</key>
  <false/>

  <key>UserName</key>
  <string>deploy</string>
  <key>WorkingDirectory</key>
  <string>/Users/deploy</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/deploy</string>
    <key>PATH</key>
    <string>{{ precis_reconcile_venv }}/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>PRECIS_PROCESS</key>
    <string>precis-reconcile</string>
    <key>PRECIS_DATABASE_URL</key>
    <string>{{ precis_reconcile_dsn }}</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/var/log/precis-reconcile.log</string>
  <key>StandardErrorPath</key>
  <string>/var/log/precis-reconcile.log</string>
</dict>
</plist>
