# deploy/roles/precis_embedder/templates/precis-embedder.service.j2
# {{ ansible_managed }}
# precis embedding service (per-node, loopback). Serves bge-m3 over HTTP
# so the same-node `precis worker` embeds via PRECIS_EMBEDDER=remote.
# Linux unit runs on the GPU node, where CUDA does the forward pass. See ADR
# 0020/0021 in the precis-mcp repo.

[Unit]
Description=precis-mcp embedding service (serve-embeddings, bge-m3)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=deploy
Group=deploy

WorkingDirectory=/tmp

Environment=PATH={{ precis_embedder_venv }}/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/var/lib/precis
Environment=HF_HOME={{ precis_embedder_hf_home }}
Environment=TOKENIZERS_PARALLELISM=false

ExecStart={{ precis_embedder_venv }}/bin/precis serve-embeddings \
  --host {{ precis_embedder_host }} \
  --port {{ precis_embedder_port }} \
  --embedder {{ precis_embedder_backend }} \
  --max-inflight {{ precis_embedder_max_inflight }}

Restart=always
RestartSec=30
StandardOutput=append:/var/log/precis-embedder.log
StandardError=append:/var/log/precis-embedder.log
# De-prioritize heavy compute so it can't starve interactive/system work
# (e.g. sshd). CPUAffinity is omitted on hosts that don't reserve cores.
Nice={{ precis_compute_nice | default(10) }}
{% if precis_compute_cpuaffinity | default('') %}CPUAffinity={{ precis_compute_cpuaffinity }}{% endif %}

[Install]
WantedBy=multi-user.target
